# bin/init-script
#! /bin/sh

### BEGIN INIT INFO
# Provides:		jarvis
# Required-Start:	$ejabberd $ircd-hybrid $syslog
# Required-Stop:	$ejabberd $ircd-hybrid $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		1
# Short-Description:	XMPP / IRC chatbots
### END INIT INFO

set -e
# /etc/init.d/jarvis: start and stop the xmpp / irc chatbots

protect_default_file(){
    if [ ! -z "${DEFAULT_FILE}" ]; then
        if [ -f "${DEFAULT_FILE}" ]; then
            if [ "$(stat -c '%U:%G' /etc/default/jarvis)" != "root:root" ]; then
                chown root:root "${DEFAULT_FILE}"
            fi
            if [ "$(stat -c '%A' /etc/default/jarvis)" != "-r--------" ]; then
                chmod 0400 "${DEFAULT_FILE}"
            fi
        fi
    fi
}


export DEFAULT_FILE="/etc/default/jarvis";
protect_default_file

if test -f ${DEFAULT_FILE}; then
    . ${DEFAULT_FILE}
fi

. /lib/lsb/init-functions

# Are we running from init?
run_by_init() {
    ([ "$previous" ] && [ "$runlevel" ]) || [ "$runlevel" = S ]
}

git_revision(){
    if [ -z ${GIT_ENABLED} ]; then 
        echo 0
        return -1
    fi
    if [ ${GIT_ENABLED} -ne 1 ]; then 
        echo 0
        return -1
    fi
    if [ ! -z "${GIT_ROOT}" -a ! -z "${GIT_REPO}" ]; then
        echo $(cd ${GIT_REPO}; git log| head -1|awk '{print $2}')
    fi
}

git_update(){
    if [ -z ${GIT_ENABLED} ]; then return -1; fi
    if [ -z ${GIT_ROOT}    ]; then return -1 ;fi
    if [ ! -d ${GIT_ROOT}  ]; then 
        mkdir -p "${GIT_ROOT}"
    fi
    if [ -z ${GIT_REPO}    ]; then return -1 ;fi
    if [ ! -d ${GIT_REPO}  ]; then 
        echo "GIT_REPO DOES NOT EXIST"
        if [ ! -z "${GIT_SOURCE}" ]; then
            (cd "${GIT_ROOT}"; git clone "${GIT_SOURCE}")
        else
            echo "GIT_SOURCE NOT DEFINED CANNOT CREATE GIT_REPO"
            return -1;
        fi
    else
        (cd ${GIT_REPO}; git pull origin master)
    fi
    # This is a temporary hack until the init script is stable:
    NOW=$(md5sum /etc/init.d/jarvis|awk '{print $2}')
    NEW=$(md5sum ${GIT_REPO}/bin/init-script|awk '{print $2}')
    if [ "${NOW}" != "${NEW}" ];then
        install -m 0744 ${GIT_REPO}/bin/init-script /etc/init.d/jarvis
    fi
}

git_rollback(){
    if [ -z ${GIT_ENABLED}]; then
        return -1
    fi
}

export PATH="${PATH:+$PATH:}/usr/sbin:/sbin"
if [ -z "${PID_FILE}" ]; then
    export PID_FILE="/var/run/jarvis.pid"
fi
if [ -z "${RUN_AS}" ]; then
    export RUN_AS='root'
fi

if [ -z "${LOG_NAME}" ]; then
    export LOG_NAME='XMPP / IRC chatbot daemon'
fi

if [ -z "${DAEMON_NAME}" ]; then
    export DAEMON_NAME='universal.pl'
fi

case "$1" in
  start)
	log_daemon_msg "Starting ${LOG_NAME}" "${DAEMON_NAME}"
        protect_default_file
        git-update
        if [ ! -f "${PID_FILE}" ]; then 
            PID_DIR=$(dirname $PID_FILE)
            if [ ! -d ${PID_DIR} ] ;then mkdir -p ${PID_DIR}; fi
            touch ${PID_FILE}
            chown ${RUN_AS} ${PID_FILE}
        fi
        COUNT=$(ps -ef | grep "${GIT_REPO}/uni[v]ersal" | awk -v runas="${RUN_AS}" '{if($1==runas){print $2}}'|wc -l)
        if [ ${COUNT} -eq 0 ]; then
	    start-stop-daemon --start --quiet --oknodo \
                              --pidfile ${PID_FILE} \
                              --chdir ${GIT_REPO} \
                              --chuid ${RUN_AS} \
                              --exec ${GIT_REPO}/universal.pl -- > /dev/null 2>&1 &
            SSD_EXIT=$?
            echo $! > ${PID_FILE}
	    if [ ${SSD_EXIT} ]; then
	        log_end_msg 0
	    else
	        log_end_msg 1
	    fi
        else
            echo "${GIT_REPO}/uni[v]ersal already running"
            exit -1
        fi
	;;
  stop)
	log_daemon_msg "Stopping ${LOG_NAME}" "${DAEMON_NAME}"
	if start-stop-daemon --stop --quiet --oknodo --pidfile ${PID_FILE}; then
	    log_end_msg 0
	else
	    log_end_msg 1
	fi
        for proc in $(ps -ef | grep "${GIT_REPO}/uni[v]ersal" | awk -v runas="${RUN_AS}" '{if($1==runas){print $2}}');do
            kill -9 ${proc}
        done
	;;

  reload|force-reload)
	log_daemon_msg "Reloading ${LOG_NAME}" "${DAEMON_NAME}"
	if start-stop-daemon --stop --signal 1 --quiet --oknodo --pidfile ${PID_FILE} --chdir ${GIT_REPO} --chuid ${RUN_AS} --exec ${GIT_REPO}/universal.pl; then
	    log_end_msg 0
	else
	    log_end_msg 1
	fi
	;;

  restart)
	log_daemon_msg "Restarting ${LOG_NAME}" "${DAEMON_NAME}"
	$0 stop
        $0 start
	;;

  status)
	status_of_proc -p ${PID_FILE} universal.pl && exit 0 || exit $?
	;;

  *)
	log_action_msg "Usage: /etc/init.d/jarvis {start|stop|reload|force-reload|restart|status}"
	exit 1
esac

exit 0
